# securitycamerapi/files/etc/cron.d/motionrestart
# Restart motion few times per day.
41 * * * * root systemctl restart motion
