# dbserver/files/etc/cron.monthly/mysqloptimize
#!/bin/sh

/usr/bin/mysqlcheck --defaults-file=/root/.my.cnf --host=localhost --optimize --auto-repair --all-databases
